# Ansible _Mysql Master Slave Replication/setup/files/supervisord.conf
[supervisord]
[program:mariadb]
command=/usr/bin/systemctl start mariadb
user=root
autostart=true
autorestart=true
