# Flask App Ansible Play book/app/files/supervisord.conf
[supervisord]
[program:fask_app]
command=/usr/bin/python /webapps/myapp/app.py
directory=/webapps/myapp/
user=myapp
autostart=true
autorestart=true
