# Osrm App Deployment Ansible code/app/files/supervisord.conf
[supervisord]
[program:osrm_app]
command=osrm-routed berlin-latest.osrm
directory=/root/osrm
user=root
autostart=true
autorestart=true
